# flaskblog.conf
[program:flaskblog]
directory=/home/yin/flaskblog
command=/home/yin/flaskblog_venv/bin/gunicorn -w 3 run:app
user=zy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

